# usr/lib/opkg/info/libopenssl-conf.list
/etc/config/openssl
/etc/ssl/openssl.cnf
/etc/init.d/openssl
